# supervisor-conf/nginx.conf
[program:nginx]
directory=/app
command=nginx -g 'daemon off;'
autostart=true
