# supervisor-conf/tor.conf
[program:tor]
command=/usr/sbin/tor -f /etc/tor/torrc
autostart=true